# containers/code/root/etc/s6-overlay/s6-rc.d/svc-openvscode-server/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# Launches openvscode-server as the abc user.

EXTENSIONS_DIR="${HOME}/extensions"
DATA_DIR="${HOME}/data"

openvscode_args=(
  --host 0.0.0.0
  --port 3000
  --without-connection-token
  --disable-telemetry
  --extensions-dir "${EXTENSIONS_DIR}"
  --user-data-dir "${DATA_DIR}"
)

if [[ -n "${VSCODE_BASE_PATH:-}" && "${VSCODE_BASE_PATH:-/}" != "/" ]]; then
  openvscode_args+=(--server-base-path "${VSCODE_BASE_PATH}")
fi

exec \
  s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z 127.0.0.1 3000" \
    cd /app/openvscode-server s6-setuidgid abc \
      /app/openvscode-server/bin/openvscode-server \
        "${openvscode_args[@]}" \
        /workspace/project
