# containers/code/root/etc/s6-overlay/s6-rc.d/init-frc-setup/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# FRC environment initialization: seeds caches and validates project mount.

GRADLE_USER_HOME="${HOME}/.gradle"
EXTENSIONS_DIR="${HOME}/extensions"
MACHINE_SETTINGS="${HOME}/data/Machine/settings.json"
SIM_LOG_FILE="${SIM_LOG_FILE:-${HOME}/sim.log}"

mkdir -p "${GRADLE_USER_HOME}" "${EXTENSIONS_DIR}" "${HOME}/data/Machine"

# Seed Gradle cache on first run.
if [[ -d /opt/frc-gradle-cache && ! -d "${GRADLE_USER_HOME}/caches" ]]; then
  echo "Seeding Gradle cache from image..."
  cp -a /opt/frc-gradle-cache/. "${GRADLE_USER_HOME}"/
fi

# Seed extensions on first run.
if [[ -d /opt/frc-extensions-cache ]] && [[ -z "$(ls -A "${EXTENSIONS_DIR}" 2>/dev/null)" ]]; then
  echo "Seeding VS Code extensions from image..."
  cp -a /opt/frc-extensions-cache/. "${EXTENSIONS_DIR}"/
fi

# Seed the remote-level editor theme. Existing non-default remote choices are preserved,
# while absent settings and the old light default move to the classroom default.
if [[ ! -s "${MACHINE_SETTINGS}" ]]; then
  printf '{\n  "workbench.colorTheme": "Dark Modern"\n}\n' > "${MACHINE_SETTINGS}"
else
  current_theme="$(jq -r '."workbench.colorTheme" // ""' "${MACHINE_SETTINGS}" 2>/dev/null || true)"
  if [[ "${current_theme}" == "" || "${current_theme}" == "Light Modern" ]]; then
    tmp_settings="$(mktemp)"
    if jq '."workbench.colorTheme" = "Dark Modern"' "${MACHINE_SETTINGS}" > "${tmp_settings}"; then
      mv "${tmp_settings}" "${MACHINE_SETTINGS}"
    else
      rm -f "${tmp_settings}"
      echo "WARNING: Could not update VS Code machine settings; leaving ${MACHINE_SETTINGS} unchanged."
    fi
  fi
fi

touch "${SIM_LOG_FILE}"

# Validate mounted project (warn only — the editor should still start).
if [[ ! -d /workspace/project ]]; then
  echo "WARNING: Mounted project directory /workspace/project does not exist."
elif [[ ! -f /workspace/project/build.gradle || ! -f /workspace/project/gradlew ]]; then
  echo "WARNING: Mounted project is missing build.gradle or gradlew."
fi

# Fix ownership for /config and /workspace (linuxserver convention).
echo "setting permissions::config"
lsiown -R abc:abc "${HOME}"
echo "setting permissions::workspace"
lsiown -R abc:abc /workspace/project
